# servers/de-fra1/etc/rc.local
#!/bin/bash

hostnamectl set-hostname de-fra1.flap42.eu
echo nameserver ::1 > /etc/resolv.conf

exit $?
