# servers/de-fra1/etc/rc.local
#!/bin/bash

hostnamectl set-hostname de-fra1.flap42.eu

exit $?
