# servers/us-nyc1/etc/rc.local
#!/bin/bash

hostnamectl set-hostname us-nyc1.flap42.eu

exit $?
