# servers/ca-yto1/etc/rc.local
#!/bin/bash

hostnamectl set-hostname ca-yto1.flap42.eu
echo nameserver ::1 > /etc/resolv.conf

exit $?
